# ss-libev/alpine/files/root/service/55-server-8989/run
#!/usr/bin/env sh

exec /usr/bin/55-server -c /config/55-server-8989/config.json 2>&1
