# ss-libev_+_unbound-tls/alpine/80/service/55-server/run
#!/usr/bin/env sh

exec /usr/bin/55-server -c /config/55-server/config.json 2>&1
