# ss-libev_+_unbound-tls/alpine/80/files/root/service/55-server/run
#!/usr/bin/env sh

exec /usr/bin/55-server -c /etc/55-server/config.json 2>&1
